# crontab


#run python script every minutes
*/1 * * * *  /usr/local/bin/python /app/project/run.py > /proc/1/fd/1 2>/proc/1/fd/2




